# deployment/equiprofile.service
[Unit]
Description=EquiProfile Node.js Application
After=network.target mysql.service mariadb.service

[Service]
Type=simple
User=www-data
WorkingDirectory=/var/equiprofile/app
EnvironmentFile=/var/equiprofile/app/.env
ExecStart=/usr/bin/node dist/index.js
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=equiprofile

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/equiprofile/app /var/log/equiprofile

[Install]
WantedBy=multi-user.target
