# deployment/ubuntu24/equiprofile.service.template
[Unit]
Description=Equiprofile Application
After=network.target

[Service]
Type=simple
User=equiprofile
WorkingDirectory=/var/www/equiprofile
Environment=NODE_ENV=production
ExecStart=/usr/bin/node dist/index.js
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
